# custom/smf/ups-nut.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type='manifest' name='site:ups-nut'>
<service name='site/ups-nut' type='service' version='1'>

        <create_default_instance enabled='true' />
        <single_instance />
        <dependency name="network" grouping="require_all" restart_on="error" type="service">
                <service_fmri value="svc:/milestone/network:default"/>
        </dependency>
        <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
                <service_fmri value='svc:/system/filesystem/local' />
        </dependency>
        <exec_method type='method' name='start' exec='/opt/custom/share/svc/ups-nut.sh start' timeout_seconds='0'>
        </exec_method>
        <exec_method type='method' name='stop' exec='/opt/custom/share/svc/ups-nut.sh stop' timeout_seconds='0'>
        </exec_method>
        <property_group name='startd' type='framework'>
                <propval name='duration' type='astring' value='transient' />
        </property_group>
        <stability value='Unstable' />
        <template>
            <common_name>
                <loctext xml:lang="C">Network UPS Tools</loctext>
            </common_name>
        </template>
</service>
</service_bundle>
